# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
User=ubuntu
Group=www-data
ListenStream=/tmp/gunicorn.sock

[Install]
WantedBy=sockets.target
